# gitserver/gitserver.conf
description "Git repository server"

start on started docker
stop on stopping docker

respawn
respawn limit 5 60

exec /usr/bin/docker run \
	-lxc-conf="lxc.network.ipv6 = 2a01:4f8:100:546f::4:30/112" \
	-lxc-conf="lxc.network.ipv6.gateway = 2a01:4f8:100:546f::4:1" \
	-v /srv/repositories:/var/git \
	zargony/gitserver
